# roles/dnsmasq/files/etc/systemd/system/dnsmasq.service.d/10-update-dhcp.conf
[Unit]
Pre-Start=/opt/ivy/set-dnsmasq-dhcp.sh
